# performance_monitor.service
[Unit]
Description=PostgreSQL Performance Monitor
Documentation=man:cron(8)
After=postgresql.service pgbouncer.service

[Service]
Type=oneshot
User=tudor
WorkingDirectory=/opt/ACTIVE/INFRA/PERFORMANCE
ExecStart=/usr/bin/python3 /opt/ACTIVE/INFRA/PERFORMANCE/enhanced_grafana_performance_metrics.py
ExecStartPost=/usr/bin/python3 /opt/ACTIVE/INFRA/PERFORMANCE/performance_alerting.py
StandardOutput=journal
StandardError=journal