# CODE/opt/ACTIVE/INFRA/BACKUPS/postgresql-backup.service
[Unit]
Description=PostgreSQL Automated Backup Service
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=tudor
Group=tudor
WorkingDirectory=/opt/ACTIVE/INFRA/BACKUPS
Environment=PYTHONPATH=/opt/ACTIVE/INFRA
Environment=TELEGRAM_BOT_TOKEN=8546618948:AAG0neoQA-kNq0M2GrZX7J-dGXNvEJEOK9w
Environment=TELEGRAM_CHAT_ID=-1002157155407
ExecStart=/usr/bin/python3 /opt/ACTIVE/INFRA/BACKUPS/postgresql_backup.py
ExecStartPost=/usr/bin/python3 /opt/ACTIVE/INFRA/BACKUPS/postgresql_rotation.py

# Resource limits
TimeoutStartSec=7200
MemoryMax=4G
CPUQuota=80%

# Restart policy
Restart=no
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/BACKUPS/postgresql /opt/ACTIVE/INFRA/LOGS /tmp

[Install]
WantedBy=multi-user.target