# BUSINESS/IDEAS/LLM/email-responder.service
[Unit]
Description=LLM Email Response Drafter
After=network.target postgresql.service

[Service]
Type=simple
User=tudor
WorkingDirectory=/opt/ACTIVE/EMAIL/LLM_RESPONDER
ExecStart=/usr/bin/python3 /opt/ACTIVE/EMAIL/LLM_RESPONDER/email_responder.py --daemon
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
