# INFRASTRUCTURE/SYSTEMD/telegram-alerts.service
[Unit]
Description=Telegram Alert Notifier for Prometheus Alertmanager
After=network.target prometheus-alertmanager.service
Wants=prometheus-alertmanager.service

[Service]
Type=simple
User=tudor
ExecStart=/usr/bin/python3 /opt/ACTIVE/INFRA/ALERTS/telegram_notifier.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-alerts

[Install]
WantedBy=multi-user.target
